# deploy/systemd/geniepod.target
[Unit]
Description=GeniePod Core Services
Documentation=https://github.com/GeniePod/genie-claw
# Audio init is optional on installs that only use chat/API paths.
Wants=genie-audio.service
After=network.target genie-audio.service

[Install]
WantedBy=multi-user.target
